# update-motd.d/colors
if [ "$NO_COLOR" -eq 1 ] ; then
    export NONE=""
    export WHITE=""
    export GREEN=""
    export RED=""
    export YELLOW=""
    export BLUE=""
    export CYAN=""
    export LIGHT_GREEN=""
    export LIGHT_RED=""
else
    export NONE="\e[m"
    export WHITE="\e[1;37m"
    export GREEN="\e[1;32m"
    export RED="\e[0;32;31m"
    export YELLOW="\e[1;33m"
    export BLUE="\e[34m"
    export CYAN="\e[36m"
    export LIGHT_GREEN="\e[1;32m"
    export LIGHT_RED="\e[1;31m"
fi
